# services/btautoconnect.service
[Unit]
Description=Bluetooth Auto Connect
After=bluetooth.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'sleep 10 && bluetoothctl connect YOUR_KEYBOARD_MAC
Restart=on-failure

[Install]
WantedBy=multi-user.target
